# scrub/xfs_scrub_fail@.service.in
# SPDX-License-Identifier: GPL-2.0-or-later
#
# Copyright (C) 2018-2024 Oracle.  All Rights Reserved.
# Author: Darrick J. Wong <djwong@kernel.org>

[Unit]
Description=Online XFS Metadata Check Failure Reporting for %f
Documentation=man:xfs_scrub(8)

[Service]
Type=oneshot
Environment=EMAIL_ADDR=root
ExecStart=@pkg_libexec_dir@/xfs_scrub_fail "${EMAIL_ADDR}" %f
User=mail
Group=mail
SupplementaryGroups=systemd-journal
